# nsqlookupd/files/conf/nsqlookupd.conf
[program:nsqlookupd]
user=root
command=nsqlookupd
autostart=true
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nsqlookupd.log
stopsignal=INT

[supervisord]
